# ansible/roles/hpc-exporter/files/hpc-exporter.service
[Unit]
Description=HPC Scheduler Exporter Service

[Service]
ExecStart=main_exporter --command 'squeue_dummy'

[Install]
WantedBy=multi-user.target
