# ansible/roles/hpc-exporter/files/hpc-exporter/hpc-exporter.service
[Unit]
Description=HPC Scheduler Exporter Service

[Service]
ExecStart=main_exporter

[Install]
WantedBy=multi-user.target
